# centos74-fludetector3/conf/supervisor.ini.j2
[program:fludetector]
command = /opt/fludetector/scripts/serve.sh
directory = /opt/fludetector
user = {{fludetector_user}}